# scripts/com.zickgraf.agentic-press.skill-metrics.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Weekly skill-usage metrics sweep for the Superpowers cherry-pick trial.

  Generated by scripts/install-skill-metrics-cron.sh from this template.
  {{WORKING_DIR}} and {{NPM_PATH}} are substituted at install time using
  the primary git worktree and `command -v npm` respectively.

  Decision date documented in src/improvements/skill-usage-report.ts:TRIAL_END_DATE.
  After the third weekly report, either:
    - extend by editing this plist's StartCalendarInterval, or
    - stop entirely:
        launchctl unload ~/Library/LaunchAgents/com.zickgraf.agentic-press.skill-metrics.plist
        rm ~/Library/LaunchAgents/com.zickgraf.agentic-press.skill-metrics.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zickgraf.agentic-press.skill-metrics</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{NPM_PATH}}</string>
        <string>run</string>
        <string>sweep-skill-metrics</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{WORKING_DIR}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{PATH}}</string>
    </dict>

    <!-- Mondays at 09:00 local time. Weekday: 1 = Monday. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Don't run on load — wait for the next scheduled time. -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{LOG_PATH}}</string>

    <key>StandardErrorPath</key>
    <string>{{LOG_PATH}}</string>
</dict>
</plist>
